# overlay/usr/lib/systemd/system/cherryWeb.service.d/env.conf
[Service]
Environment="PERL5LIB=/var/lib/cherryepg/cherryTool/lib:/var/lib/cherryepg/perl5/lib/perl5"
Environment="DANCER_ENVIRONMENT=production"
