# overlay/usr/lib/systemd/system/cherryWeb.service.d/env.conf
[Service]
Environment="PATH=/var/lib/cherryepg/perl5/bin:/usr/local/bin:/usr/bin:/bin:/usr/local/games:/usr/games:/var/lib/cherryepg/cherryTool/bin"
Environment="PERL5LIB=/var/lib/cherryepg/perl5/lib/perl5:/var/lib/cherryepg/cherryTool/lib"
Environment="PERL_LOCAL_LIB_ROOT=/var/lib/cherryepg/perl5"
Environment="DANCER_ENVIRONMENT=production"
